# debian/postinst
if ! getent passwd passwall >/dev/null; then
        adduser --system --group --force-badname --quiet \
            --home /run/passwall/ --no-create-home \
            --shell /bin/false \
            passwall
fi

ln -s /usr/share/applications/passwall.desktop /etc/xdg/autostart/passwall.desktop
cp /usr/share/passwall/store/config-sample.yml /usr/share/passwall/store/config.yml
mkdir -p /etc/passwall
ln -s /usr/share/passwall/store/config.yml /etc/passwall/config.yml

chown -R passwall:passwall /usr/share/passwall/
chmod -R 770 /usr/share/passwall/

echo """
[Unit]
Description=Passwall Server
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=passwall
Group=passwall
ExecStart=/usr/bin/passwall/server/passwall-server
Environment=\"GIN_MODE=release\"
WorkingDirectory=/usr/share/passwall/
[Install]
WantedBy=multi-user.target
    """ > /etc/systemd/system/passwall-server.service

systemctl daemon-reload
systemctl enable passwall-server
systemctl restart passwall-server


echo "Your settings are stored in /etc/passwall/config.yml"

printf "You can login with: \nUsername: passwall \nPassword: password\n"

printf "Web address: http://localhost:3625 \n"